# .config/systemd/user/dns.service
[Unit]
Description=update dns record

[Service]
Type=oneshot
ExecStart=/home/evan/bin/dns

[Install]
WantedBy=suspend.target