# .config/systemd/user/redrum.timer
[Unit]
Description=Run imgurt.service periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
